# etc/systemd/network-check.timer
[Unit]
Description=Run network-check periodically

[Timer]
OnBootSec=5min
OnUnitActiveSec=60

[Install]
WantedBy=timers.target
